# deploy/a1/keepalive.timer
[Unit]
Description=Run keep-alive CPU burst every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s

[Install]
WantedBy=timers.target
