# deploy/a1/llm-healthcheck.timer
[Unit]
Description=Run LLM health check every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
